# cronjob.cron
0 * * * * python proxypaper.py && python proxy.py
